# configs/systemd/user/clipmenud.service.d/override.conf
[Service]
Environment=CM_SELECTIONS=clipboard
Environment=CM_IGNORE_WINDOW=KeePassXC$
Environment=CM_DIR=.cache/clipmenu
