# configs/systemd/system/betterlockscreen@.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/betterlockscreen -l dim -t "Don't touch my machine!"
